# dow_serv2.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn
script
chdir /root/DoW/dawnofworlds
exec python3 cherrypy_attemp.py
end script